# supervisord.conf
; Runs php-fpm + apache + the Nextcloud log tailer in one container.
; Either dying gets restarted; previously apache ran alone via the compose
; `command:` and php was in-process (mod_php).
[supervisord]
nodaemon=true
user=root
pidfile=/run/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

; /dev/fd/1 + maxbytes=0 routes each program's output to the container log
; (supervisord's own stdout) without the seek errors /dev/stdout triggers.

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=10

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=20

[program:nextcloud-logs]
command=/usr/bin/tail -F /var/www/html/data/nextcloud.log /var/www/html/data/audit.log
autorestart=true
startsecs=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=30
